# deployment/raft-node.service
[Unit]
Description=Raft Control Plane Node
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/raft-control-plane
ExecStart=/opt/raft-control-plane/bin/raft-node --config /etc/raft/config.json
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
